# deploy/sukesan.service
# systemd ユニットの例（Linux）。
#
# 使い方:
#   1. __APP_DIR__ をアプリの絶対パス、__USER__ を実行ユーザー名に置き換える
#   2. sudo cp deploy/sukesan.service /etc/systemd/system/sukesan.service
#   3. sudo systemctl daemon-reload
#   4. sudo systemctl enable --now sukesan
#   状態確認: systemctl status sukesan / ログ: journalctl -u sukesan -f
#
# 注意: rbenv / rvm 等を使う場合、ruby/bundle が PATH に無いと起動に失敗する。
#       その場合は Environment=PATH=... を追記するか、ExecStart を絶対パスの bundle に変更する。

[Unit]
Description=SUKESAN scheduling tool
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__APP_DIR__
Environment=APP_ENV=production
# .env は dotenv によりアプリ起動時に読み込まれる
ExecStart=__APP_DIR__/bin/server run
Restart=on-failure
RestartSec=3

# --- サンドボックス（最小権限化）---
# アプリは data/ log/ tmp/ 以外への書き込みが不要なので、ファイルシステムを読み取り専用に寄せる。
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# ProtectHome=true は /home 配下を不可視にするため、__APP_DIR__ は /opt/sukesan など /home の外に置くこと
# （/home 配下に置くとアプリが自身のファイルを読めず起動に失敗する）。/home 配下で運用する場合は
# ProtectHome=read-only に緩めるか、この行を外す。
ProtectHome=true
# 書き込みを許可するのはアプリのデータ/ログ/PID 用ディレクトリのみ（PID は tmp/pids/ に出る）。
ReadWritePaths=__APP_DIR__/data __APP_DIR__/log __APP_DIR__/tmp
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true
SystemCallArchitectures=native
# 実行可能メモリへの書き込みを禁止。MRI（インタプリタ実行）では問題ないが、YJIT/JIT を有効化する場合は
# 実行コード生成と競合するため、この行を外すこと。
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
